# docker-file/docker-centos-supervisor-master/supervisor.d/dovecot.conf
[program:dovecot]
command = /usr/sbin/dovecot
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE

