# docker-file/docker-centos-supervisor-master/supervisor.d/postfix.conf
[program:postfix]
command = /usr/sbin/postfix start
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE

