# docker-file/docker-centos-supervisor-master/supervisor.d/rsyslog.conf
[program:rsyslog]
command = /usr/sbin/rsyslogd -n
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE

